# bsp/arch-dm368/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttyS1::askfirst:-/bin/sh
