# bsp/mach-leopard/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttyS0::askfirst:-/bin/sh
